# debian/pmg-pbsbackup@.service
[Unit]
Description=Backup to PBS remote %I

[Service]
Type=oneshot
ExecStart=/usr/bin/pmgbackup proxmox-backup backup %I
